# ansible/galaxy_roles/erik-moroz.Supervisor/templates/program.conf.j2
; {{ ansible_managed }}

[program:{{ item.value.name }}]
{% for option, value in item.value.items() %}
{% if option != 'name' %}
{{ option }} = {{ value }}
{% endif %}
{% endfor %}
